# files/chancoin-bootnode.service
[Unit]
Description=Chancoin go client

[Service]
User=varg
Group=varg
Type=simple
ExecStart=/opt/chancoin/build/bin/bootnode -nodekey ./bootnode.key

[Install]
WantedBy=default.target
